# extra-coreos/usr/lib/systemd/system/systemd-fsck@dev-disk-by\x2dlabel-boot.service.d/override.conf
[Unit]
After=boot-part-discovered.target
ConditionFileNotEmpty=/run/coreos/boot-partition-status
